# data/supercronic/crontab
30 * * * *   /usr/local/bin/ampache_cron.sh
